# panasonic-kb.service
[Unit]
Description=Panasonic keyboard thing

[Service]
ExecStart=/path/to/panasonic-kb

[Install]
WantedBy=multi-user.target